# scripts/beep-analytics.service
[Unit]
Description=beep-analytics server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=beep-analytics
Group=beep-analytics
ExecStart=/usr/bin/beep-analytics serve \
  --port 8080 \
  --db /var/lib/beep-analytics/beep.db \
  --geoip /usr/share/GeoIP
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/beep-analytics
ReadOnlyPaths=/usr/share/GeoIP

[Install]
WantedBy=multi-user.target
